# templates/process.conf.erb
start on starting <%= application %>
stop  on stopping <%= application %>

respawn

setuid <%= user %>
setgid <%= user %>

env HOME=<%= home %>
<%- environment_variables.each do |environment_variable| %>
env <%= environment_variable %>
<% end -%>

chdir <%= application_root %>

exec bash -lc 'chruby-exec ruby -- <%= process.command %> >> <%= log %>/<%= application %>/<%= process.name %>.log 2>&1'
